# build/config/includes.chroot/etc/systemd/system/unattended-upgrades.service.d/shiros-skip-on-live.conf
[Unit]
# shirOS — no correr unattended-upgrades en sesión live.
ConditionKernelCommandLine=!boot=live
