# build/config/includes.chroot/etc/systemd/system/lightdm.service.d/shiros-skip-on-live.conf
[Unit]
# shirOS — no levantar LightDM en la sesión live.
# En live boot usamos getty@tty1 autologin + startxfce4 directo.
# En el sistema instalado (sin boot=live) esta condición es falsa,
# así que el drop-in no aplica y LightDM arranca normal.
ConditionKernelCommandLine=!boot=live
