# build/config/includes.chroot/etc/systemd/system/getty@tty3.service.d/shiros-rescue.conf
[Unit]
# shirOS — TTY3 de rescate.
# SIEMPRE corre como root con autologin, sin contraseña. Pensado para
# cuando el desktop o LightDM rompen y necesitás shell para debug.
# Solo aplica en sesión live (no en sistema instalado).
ConditionKernelCommandLine=boot=live

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
