# build/config/includes.chroot/etc/systemd/system/apt-daily-upgrade.timer.d/shiros-skip-on-live.conf
[Unit]
# shirOS — no correr apt-daily-upgrade en sesión live.
ConditionKernelCommandLine=!boot=live
