# build/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/shiros-live-autologin.conf
[Unit]
# shirOS — autologin como 'live' en tty1, solo durante boot live.
# Si no hay boot=live en cmdline (= sistema instalado), este drop-in
# no aplica y getty pide login normalmente.
ConditionKernelCommandLine=boot=live
After=live-config.service shiros-live-init.service

[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \\u' --autologin live --noclear --keep-baud %I 115200,38400,9600 $TERM
