# build/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/shiros-live-autologin.conf
[Unit]
# shirOS — autologin como 'live' en tty1, solo durante boot live.
# En sistema instalado (sin boot=live) el drop-in no aplica.
ConditionKernelCommandLine=boot=live

[Service]
# Forma minimal de autologin. NO usar `-o '-p -- \u'`: ese argumento agrega
# `-- live` al cmdline de login, lo cual termina el parseo de opciones y
# anula el `-f live` que agetty agrega por --autologin. Resultado: login
# pide password normal en vez de saltearlo. Bug visto en v0.1.4.
ExecStart=
ExecStart=-/sbin/agetty --autologin live --noclear %I $TERM
