# build/config/includes.chroot/etc/systemd/system/apt-daily.timer.d/shiros-skip-on-live.conf
[Unit]
# shirOS — no correr apt-daily en sesión live (es efímera).
# En sistema instalado (sin boot=live) el drop-in no aplica → corre normal.
ConditionKernelCommandLine=!boot=live
